# development/grunt.conf
[program:grunt]
directory=/vagrant/
command=/vagrant/development/grunt.sh
autostart=true
autorestart=true